# sqs2mandrill/templates/default/sqs2mandrill-sqs2mandrill.conf
start on starting sqs2mandrill
stop on stopping sqs2mandrill
