# sqs2cb/templates/default/sqs2cb-sqs2cb.conf
start on starting sqs2cb
stop on stopping sqs2cb
